# python3-ros/kinetic/runable/.private/rosbash/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/home/parallels/Desktop/compare2/kinetic/src/ros/rosbash/roszsh"
