# python3-ros/kinetic/runable/.private/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/parallels/Desktop/compare2/kinetic/src/ros/rosbash/rostcsh"
