# python3-ros/kinetic/runable/etc/catkin/profile.d/1.ros_package_path.sh
/home/parallels/Desktop/compare2/kinetic/devel/.private/ros_environment/etc/catkin/profile.d/1.ros_package_path.sh