# python3-ros/kinetic/runable/etc/catkin/profile.d/10.rosconsole.sh
/home/parallels/Desktop/compare2/kinetic/devel/.private/rosconsole/etc/catkin/profile.d/10.rosconsole.sh