# python3-ros/kinetic/runable/etc/catkin/profile.d/10.rosbuild.sh
/home/parallels/Desktop/compare2/kinetic/devel/.private/rosbuild/etc/catkin/profile.d/10.rosbuild.sh