# python3-ros/kinetic/runable/etc/catkin/profile.d/10.roslaunch.sh
/home/parallels/Desktop/compare2/kinetic/devel/.private/roslaunch/etc/catkin/profile.d/10.roslaunch.sh