# python3-ros/kinetic/runable/etc/catkin/profile.d/15.rosbash.bash
/home/parallels/Desktop/compare2/kinetic/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.bash