# python3-ros/kinetic/runable/etc/catkin/profile.d/20.transform.bash
/home/parallels/Desktop/compare2/kinetic/devel/.private/topic_tools/etc/catkin/profile.d/20.transform.bash