# python3-ros/kinetic/runable/etc/catkin/profile.d/99.roslisp.sh
/home/parallels/Desktop/compare2/kinetic/devel/.private/roslisp/etc/catkin/profile.d/99.roslisp.sh