# python3-ros/kinetic/runable/etc/catkin/profile.d/cmake.lock
/home/parallels/Desktop/compare2/kinetic/devel/.private/ros_environment/etc/catkin/profile.d/cmake.lock